# dist/rpmbuild.target
[Unit]
Description=Telego telegram MTPRoto service
